# ansible/roles/sanoid/templates/syncoid-replication.service.j2
[Unit]
Description=Syncoid ZFS Replication
After=zfs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/syncoid-replication.sh
